# supervisord.conf
[supervisord]
nodaemon=true

[program:folderserve]
command=/serve.sh
autorestart=true

[program:jdownloader]
command=/jd.sh
autorestart=true
